# .config/systemd/user/clipcat.service
[Unit]
Description=Clipcat Daemon
PartOf=default.target

[Install]
WantedBy=default.target

[Service]
ExecStartPre=/bin/rm -f %t/clipcat/grpc.sock
ExecStart=/usr/bin/clipcatd --no-daemon --replace
Restart=on-failure
Type=simple
