# .config/systemd/user/feh.timer
[Unit]
Description=Run feh.service every 15 minutes

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
