# .config/systemd/user/low_battery.service
[Unit]
Description=Notifications for low battery

[Service]
Type=oneshot
ExecStart=/home/sugimoto/.local/bin/low-battery
